# conf/liagent.sv.conf
[program:agent]
command=/bin/bash -c "service liagentd restart"
redirect_stderr=true
